# src/init.d/ssh-remote
#!/sbin/openrc-run
# Copyright 2025 Andrei Ruslantsev
# Distributed under the terms of the GNU General Public License v3

supervisor="supervise-daemon"

depend() {
	use logger dns net sshd
	CFG_NAME=$(echo ${RC_SVCNAME} | cut -d '.' -f 2)
	command="/usr/bin/ssh-remote /etc/conf.d/ssh-remote.${CFG_NAME}"
}
